# deb/etc/profile.d/golang.sh
export GOTMPDIR="/usr/tmp"
export PATH="$PATH:/usr/local/libexec/go/bin"